# PBL_FINALPROJECT/f12b91b8-9bfb-49dc-8725-6d54a62c8ac8/project-files/docker/ac6218f3-675d-4035-aaf4-e84d97b6ef4a/etc/update-motd.d/00-header
version https://git-lfs.github.com/spec/v1
oid sha256:540acb8fe88212157b1d91348d96e2d514b079ca9559fbb32bc23e604676f06e
size 1249
