# PBL_FINALPROJECT/f12b91b8-9bfb-49dc-8725-6d54a62c8ac8/project-files/docker/ac6218f3-675d-4035-aaf4-e84d97b6ef4a/etc/update-motd.d/10-help-text
version https://git-lfs.github.com/spec/v1
oid sha256:f4a274710c126c70e4b09ddd1a042965fc3304b4acd699b32048d902969b61bb
size 1183
