# PBL_FINALPROJECT/f12b91b8-9bfb-49dc-8725-6d54a62c8ac8/project-files/docker/ac6218f3-675d-4035-aaf4-e84d97b6ef4a/etc/update-motd.d/50-motd-news
version https://git-lfs.github.com/spec/v1
oid sha256:00fa958650aadc460c639feca8fc89dfe13b6e83661a8caa622ad94f6cc0de73
size 5169
