# PBL_FINALPROJECT/f12b91b8-9bfb-49dc-8725-6d54a62c8ac8/project-files/docker/ac6218f3-675d-4035-aaf4-e84d97b6ef4a/etc/update-motd.d/60-unminimize
version https://git-lfs.github.com/spec/v1
oid sha256:c697931dc7b753ef83fa5f7092c20d7ccf3f68841b384cc1d2d556b1b7f85d27
size 365
